# deploy/mt5-wine.service
[Unit]
Description=MetaTrader 5 (headless via Wine + Xvfb)
After=network-online.target ftmo-monitor.service
Wants=network-online.target

[Service]
Type=simple
User=root
Environment=DISPLAY=:99
Environment=WINEPREFIX=/root/.wine_mt5
Environment=WINEARCH=win64
Environment=WINEDEBUG=-all

# Start virtual display, then MT5
ExecStartPre=/bin/bash -c "pkill Xvfb || true; sleep 1"
ExecStartPre=/usr/bin/Xvfb :99 -screen 0 1024x768x24 -ac +extension GLX +render &
ExecStartPre=/bin/sleep 3

ExecStart=/usr/bin/wine "/root/.wine_mt5/drive_c/Program Files/MetaTrader 5/terminal64.exe" /portable

Restart=always
RestartSec=15
TimeoutStopSec=30

# Log to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mt5-wine

[Install]
WantedBy=multi-user.target
